# src/sunrise-alarm.service
[Unit]
Description=sunlight-alarm: A nice alarm clock

[Service]
WorkingDirectory=/opt/sunrise-alarm/app/
ExecStart=/usr/bin/node /opt/sunrise-alarm/app/sunlight-alarm
Restart=always
           
[Install]
WantedBy=default.target
